# systemd/pihole-cputemp-json.service
[Unit]
Description=Write Pi-hole CPU temperature JSON

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/write-pihole-cputemp-json
